# pkg/debian/yeti_clear_obsolete_dumps.cron
# /etc/cron.d/yeti-clear-obsolete-dumps : crontab fragment for yeti-clear-obsolete-dumps

# checking for event in database every hour
0 * * * * root [ -x /usr/sbin/yeti-remove-obsolete-dumps ] && /usr/sbin/yeti-remove-obsolete-dumps clean
